# crontab
# Указываем путь к python внутри python:3.11-slim
PATH=/usr/local/bin:/usr/bin:/bin

# Обработка сообщений Telegram каждые 30 минут
*/30 * * * * cd /app && python telegram_to_sheets.py >> /app/logs/telegram_to_sheets.log 2>&1

# Загрузка кода из Bitbucket каждый день в 6:00
0 6 * * * cd /app && echo "---- $(date '+%Y-%m-%d %H:%M:%S') fetch_code_from_bitbucket.py ----" >> /app/logs/fetch_code_cron.log && python fetch_code_from_bitbucket.py >> /app/logs/fetch_code_cron.log 2>&1

# Обработка необработанных ошибок через GPT каждый день в 6:10  
10 6 * * * cd /app && echo "---- $(date '+%Y-%m-%d %H:%M:%S') process_unhandled_errors.py ----" >> /app/logs/gpt_process.log && python process_unhandled_errors.py >> /app/logs/gpt_process.log 2>&1

# Отправка ежедневного отчета каждый день в 6:19
19 6 * * * cd /app && echo "---- $(date '+%Y-%m-%d %H:%M:%S') unknown_transaction.py ----" >> /app/logs/unknown_tx.log && python unknown_transaction.py >> /app/logs/unknown_tx.log 2>&1

# Отправка ежедневного отчета каждый день в 6:20
20 6 * * * cd /app && echo "---- $(date '+%Y-%m-%d %H:%M:%S') send_daily_summary.py ----" >> /app/logs/daily_summary.log && python send_daily_summary.py >> /app/logs/daily_summary.log 2>&1

# Пустая строка в конце файла (обязательна для cron)
